# tproxy.service
[Service]
Type=oneshot
EnvironmentFile=/opt/xray/.env
ExecStart=/opt/xray/tproxy.sh
RemainAfterExit=yes
[Install]
WantedBy=multi-user.target
