# tproxy.service
[Service]
Type=oneshot
EnvironmentFile=/opt/xray/.env
ExecStart=/usr/local/tproxy.sh
RemainAfterExit=yes
[Install]
WantedBy=multi-user.target
